# spec/resources/export/upstart/app-bravo-1-custom.conf
start on starting app-bravo
stop on stopping app-bravo
respawn

env PORT=5100

setuid app

reload signal HUP

kill signal QUIT

chdir /tmp/app

exec ./bravo
