# spec/resources/export/upstart/app-alpha-1.conf
start on starting app-alpha
stop on stopping app-alpha
respawn

env PORT=5000

setuid app

reload signal HUP

kill signal TERM

chdir /tmp/app

exec ./alpha
